# automountnfs-check.timer
[Unit]
Description=Checks NFS mounts every minute.

[Timer]
OnCalendar=minutely

[Install]
WantedBy=timers.target automountnfs.service
